# backend/deployment/systemd/kavan-celery.service
[Unit]
Description=KAVAN v6.0 — Celery Worker
After=network.target redis.service postgresql.service
Wants=redis.service postgresql.service

[Service]
Type=forking
User=kavan
Group=kavan
WorkingDirectory=/opt/kavan/backend
EnvironmentFile=/opt/kavan/.env
ExecStart=/opt/kavan/venv/bin/celery \
    -A config.celery \
    multi start 2 \
    --pidfile=/run/kavan/celery-%n.pid \
    --logfile=/var/log/kavan/celery-%n.log \
    --loglevel=INFO \
    --concurrency=4
ExecStop=/opt/kavan/venv/bin/celery \
    -A config.celery \
    multi stopwait 2 \
    --pidfile=/run/kavan/celery-%n.pid \
    --logfile=/var/log/kavan/celery-%n.log
ExecReload=/opt/kavan/venv/bin/celery \
    -A config.celery \
    multi restart 2 \
    --pidfile=/run/kavan/celery-%n.pid \
    --logfile=/var/log/kavan/celery-%n.log \
    --loglevel=INFO
RuntimeDirectory=kavan
RuntimeDirectoryMode=0755

StandardOutput=journal
StandardError=journal
SyslogIdentifier=kavan-celery

Restart=on-failure
RestartSec=10s

[Install]
WantedBy=multi-user.target
